# gocd-server/supervisor/gocd-server.conf
[program:gocd-server]
# https://github.com/sameersbn/docker-squid/blob/master/start
command=java \
          -XX:MaxPermSize=256m \
          -Xmx1024m \
          -Dserver.hostname=$RUNDECK_HOSTNAME \
          -Ddefault.user.password=$RUNDECK_USER_PASSWORD \
          -Drundeck.jetty.connector.forwarded=$RUNDECK_FORWARDED \
          -jar $RDECK_BASE/rundeck-launcher.jar --skipinstall \
          -c /prj/data/rundeck/config \
          -p /prj/data/rundeck/projects
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
startsecs=10
stopwaitsecs=10
