# rundeck/supervisor/rundeck.conf
[program:rundeck]
# https://github.com/sameersbn/docker-squid/blob/master/start
command=/usr/bin/java
    -XX:MaxPermSize=256m
    -Xmx1024m
    -Dserver.hostname=%(ENV_RUNDECK_HOSTNAME)s
    -Drundeck.jetty.connector.forwarded=%(ENV_RUNDECK_FORWARDED)s
    -Dserver.web.context=%(ENV_RUNDECK_START_URL)s
    -jar %(ENV_RDECK_BASE)s/rundeck-launcher.jar --skipinstall
    -c /prj/data/rundeck/config
    -p /prj/data/rundeck/projects
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
startsecs=10
stopwaitsecs=10
