# monitor-icinga2/supervisor/icinga2.conf
[program:icinga2]
command=icinga2 daemon
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
startretries=0
