# proxy-cache/supervisor/squid.conf
[program:squid]
# https://github.com/sameersbn/docker-squid/blob/master/start
command=/usr/sbin/squid3 -NYC -d 1 -f /etc/squid3/squid.conf
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
startsecs=10
stopwaitsecs=10
